# v3.1/deploy/spider-brain.service
[Unit]
Description=Spider Robot v3.1 Brain Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/root
ExecStart=/root/brain_daemon
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spider-brain

[Install]
WantedBy=multi-user.target
